# service/hiddify
#!/bin/sh /etc/rc.common

START=99
STOP=10

RAM_MODE=true

HIDDIFY_URL="https://github.com/AmirParsaRabiei/hiddify-core/releases/download/v1.0.1/hiddify-cli-linux-mipsel-softfloat.tar.gz"
# HIDDIFY_URL="https://github.com/PacketCipher/hiddify-core/releases/download/draft/hiddify-cli-linux-mipsel-softfloat.tar.gz"
# HIDDIFY_URL="https://github.com/hiddify/hiddify-core/releases/download/v3.1.8/hiddify-cli-linux-mipsel-softfloat.tar.gz"

SUB_URL="https://raw.githubusercontent.com/AmirParsaRabiei/TVC/refs/heads/main/subscriptions/xray/normal/mix"
# SUB_URL="https://raw.githubusercontent.com/hiddify/hiddify-next/main/test.configs/warp"
# SUB_URL="https://raw.githubusercontent.com/PacketCipher/TVC/main/subscriptions/hiddify/warp"
# SUB_URL="https://raw.githubusercontent.com/hiddify/hiddify-app/refs/heads/main/test.configs/warp"

# Modify OpenVPN to be enabled or not
Enable_OVPN=false # This script part focuses on non-OVPN part. OVPN watchdog would need similar changes.

# CRON_JOB_OLD="18,48 * * * * /etc/init.d/hiddify restart" # Old cron job
CRON_JOB_NEW="18,48 * * * * /etc/init.d/hiddify update_config" # New cron job

LIVE_CONFIG_FILE="/root/hiddify-sub"
WATCHDOG_SCRIPT_NAME="hiddify_watchdog.sh" # Assuming this is the non-OpenVPN watchdog

if [ "$RAM_MODE" = true ]; then
    SERVICE_DIR="/tmp/usr/bin/HiddifyCli"
else
    SERVICE_DIR="/root/HiddifyCli"
fi

download_and_extract() {
    echo "Downloading and extracting Hiddify CLI..."
    [ ! -d "$SERVICE_DIR" ] && mkdir -p "$SERVICE_DIR"
    local retries=10
    local count=0
    until [ $count -ge $retries ]
    do
        wget -O - "$HIDDIFY_URL" | tar -xz -C "$SERVICE_DIR" && break
        count=$((count+1))
        echo "Retrying... ($count/$retries)"
        sleep 60
    done
    if [ $count -ge $retries ]; then
        echo "Failed to download and extract Hiddify CLI after $retries attempts."
        return 1
    fi
}

# Function to download initial subscription config if it doesn't exist
download_initial_config() {
    if [ ! -f "$LIVE_CONFIG_FILE" ]; then
        echo "Initial configuration $LIVE_CONFIG_FILE not found. Downloading..."
        curl -L -o "$LIVE_CONFIG_FILE" "$SUB_URL"
        if [ ! -s "$LIVE_CONFIG_FILE" ]; then
            echo "Error: Downloaded initial configuration file $LIVE_CONFIG_FILE is empty or download failed."
            rm -f "$LIVE_CONFIG_FILE" # Clean up
            return 1
        else
            echo "Initial configuration downloaded successfully."
        fi
    else
        echo "Configuration file $LIVE_CONFIG_FILE already exists."
    fi
    return 0
}

add_cron_job() {
    echo "Adding/Updating cron job to: $CRON_JOB_NEW"
    # Remove old cron job if it exists, then add new one
    (crontab -l 2>/dev/null | grep -v -F "/etc/init.d/hiddify restart" | grep -v -F "$CRON_JOB_NEW"; echo "$CRON_JOB_NEW") | crontab -
}

remove_cron_job() {
    echo "Removing cron job..."
    crontab -l 2>/dev/null | grep -v -F "$CRON_JOB_NEW" | crontab -
}

start() {
    echo "Starting Hiddify service..."
    
    if [ ! -f "$SERVICE_DIR/HiddifyCli" ]; then
        download_and_extract || return 1
    fi

    # Download initial config if it's not present
    download_initial_config || return 1

    echo "Starting $WATCHDOG_SCRIPT_NAME..."
    # Pass SUB_URL and SERVICE_DIR as environment variables to the watchdog
    # The watchdog script has been modified to use these from the environment
    if [ "$Enable_OVPN" = true ]; then
        # This part is for OVPN, which is out of scope for current changes,
        # but if it were to be updated, hiddify_openvpn_watchdog.sh would need similar signal handling.
        echo "OpenVPN mode is selected. Ensure hiddify_openvpn_watchdog.sh is adapted for new update logic."
        SUB_URL=$SUB_URL SERVICE_DIR=$SERVICE_DIR nohup /root/hiddify_openvpn_watchdog.sh > /dev/null 2>&1 &
    else
        SUB_URL=$SUB_URL SERVICE_DIR=$SERVICE_DIR nohup "/root/$WATCHDOG_SCRIPT_NAME" > /dev/null 2>&1 &
    fi

    add_cron_job
}

stop() {
    echo "Stopping Hiddify service..."
    if [ "$Enable_OVPN" = true ]; then
        echo "Stopping OpenVPN related processes..."
        pkill -f hiddify_openvpn_watchdog.sh # Use pkill for safety
        killall HiddifyCli # Assuming HiddifyCli is common
        killall openvpn
    else
        echo "Stopping standard Hiddify processes..."
        # Use pkill -f to ensure the correct watchdog script is killed.
        # This is safer if other scripts might have similar names.
        pkill -f "/root/$WATCHDOG_SCRIPT_NAME"
        # It's also good practice to ensure HiddifyCli is stopped directly,
        # as watchdog might be killed before it stops HiddifyCli.
        pkill -f "HiddifyCli run --config $LIVE_CONFIG_FILE"
    fi
    remove_cron_job
}

restart() {
    echo "Restarting Hiddify service..."
    stop
    sleep 1
    start
}

update_config() {
    echo "Attempting to trigger configuration update for Hiddify..."
    local watchdog_pid
    # Find the PID of the correct watchdog script
    # We use pgrep -f to match the full path, avoiding issues with other scripts named similarly
    watchdog_pid=$(pgrep -f "/root/$WATCHDOG_SCRIPT_NAME")

    if [ -n "$watchdog_pid" ]; then
        # Send SIGUSR1 signal to the watchdog script
        # Handle cases where multiple PIDs might be found (though unlikely for nohup'd script)
        for pid in $watchdog_pid; do
            echo "Sending SIGUSR1 to $WATCHDOG_SCRIPT_NAME (PID: $pid)"
            kill -USR1 "$pid"
        done
        echo "Configuration update signal sent."
    else
        echo "Error: $WATCHDOG_SCRIPT_NAME is not running. Cannot send update signal."
        echo "Consider starting the service."
    fi
}

# Standard OpenWrt service functions
start_service() {
    start
}

stop_service() {
    stop
}

disable() {
    stop
    # Optionally, could also remove the HiddifyCli binaries/configs if disable means full removal
    # For now, just stopping and removing cron job.
    echo "Hiddify service disabled."
}

# Note: For this new system to work, the cron job needs to be updated from
# "18,48 * * * * /etc/init.d/hiddify restart"
# to
# "18,48 * * * * /etc/init.d/hiddify update_config"
# The add_cron_job function has been updated to reflect this.
