# service/hiddify
#!/bin/sh /etc/rc.common

START=99
STOP=10

SERVICE_DIR="/tmp/usr/bin"
HIDDIFY_URL="https://github.com/hiddify/hiddify-core/releases/download/v1.2.1/hiddify-cli-linux-mipsel-softfloat.tar.gz"

# Modify OVPN to be enabled or not
Enable_OVPN=true

CRON_JOB="0 * * * * /etc/init.d/hiddify restart"

download_and_extract() {
    echo "Downloading and extracting Hiddify CLI..."

    # Create directory if it doesn't exist
    [ ! -d "$SERVICE_DIR" ] && mkdir -p "$SERVICE_DIR"

    # Download and extract the Hiddify CLI with retry logic
    local retries=10
    local count=0
    until [ $count -ge $retries ]
    do
        wget -O - "$HIDDIFY_URL" | tar -xz -C "$SERVICE_DIR" && break
        count=$((count+1))
        echo "Retrying... ($count/$retries)"
        sleep 60
    done

    if [ $count -ge $retries ]; then
        echo "Failed to download and extract Hiddify CLI after $retries attempts."
        return 1
    fi
}

add_cron_job() {
    echo "Adding cron job..."
    (crontab -l 2>/dev/null | grep -F "$CRON_JOB") || (crontab -l 2>/dev/null; echo "$CRON_JOB") | crontab -
}

remove_cron_job() {
    echo "Removing cron job..."
    crontab -l 2>/dev/null | grep -v -F "$CRON_JOB" | crontab -
}

start() {
    echo "Starting Hiddify service..."
    
    # Check if the HiddifyCli exists, if not, download and extract
    if [ ! -f "$SERVICE_DIR/HiddifyCli" ]; then
        download_and_extract || return 1
    fi

    # Run the watchdog script with nohup
    if [ "$Enable_OVPN" = true ]; then
        nohup /root/hiddify_openvpn_watchdog.sh > /dev/null 2>&1 &
    else
        nohup /root/hiddify_watchdog.sh > /dev/null 2>&1 &
    fi

    add_cron_job
}

stop() {
    echo "Stopping Hiddify service..."
    if [ "$Enable_OVPN" = true ]; then
        killall hiddify_openvpn_watchdog.sh
        killall HiddifyCli
        /etc/init.d/openvpn stop
    else
        killall hiddify_watchdog.sh
        killall HiddifyCli
    fi

    remove_cron_job
}

restart() {
    stop
    sleep 1
    start
}

start_service() {
    start
}

stop_service() {
    stop
}

disable() {
    stop
}
